# etc/systemd/iptsd.service.in
[Unit]
Description=Intel Precise Touch & Stylus Daemon
Documentation=https://github.com/linux-surface/iptsd
After=dev-ipts-15.device
Wants=dev-ipts-15.device

[Service]
Type=simple
ExecStart=@bindir@/iptsd

[Install]
WantedBy=multi-user.target
