# etc/systemd/iptsd@.service.in
[Unit]
Description=Intel Precise Touch & Stylus Daemon
Documentation=https://github.com/linux-surface/iptsd

[Service]
Type=simple
ExecStart=@bindir@/iptsd %I
